# base/root/etc/services.d/openssh-server/run
#!/usr/bin/with-contenv bash

USERNAME=${USERNAME:-"eyhn"}

exec 2>&1 \
        s6-setuidgid ${USERNAME} /usr/sbin/sshd -D -e -p 2222